# streamdeck/systemd/foodassistant-streamdeck.service
[Unit]
Description=FoodAssistant Stream Deck controller
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
# Run as a normal user that belongs to the plugdev group (see the udev rule)
# so the controller can claim the Stream Deck without root.
User=foodassistant
Group=plugdev
WorkingDirectory=/opt/foodassistant
Environment=FOODASSISTANT_BASE_URL=http://localhost:9284
# Keep the API key out of the repo and the config file. Point this at a file
# holding a single line: FOODASSISTANT_API_KEY=...
EnvironmentFile=-/etc/foodassistant/streamdeck.env
ExecStart=/opt/foodassistant/venv/bin/python -m foodassistant_streamdeck --config /etc/foodassistant/streamdeck.toml
# on-failure (not always) on purpose: the controller recovers a lost or
# unplugged deck in-process (it idle-waits for the deck to return and
# re-initialises it), so systemd is only the backstop for a hard crash. A
# missing deck at startup is no longer a nonzero exit, so this policy does not
# thrash-restart while a deck is absent; in-process recovery stays the primary
# mechanism (FoodAssistant-o29k).
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
